# gethosts
#!/bin/sh /etc/rc.common

START=99
STOP=10
USE_PROCD=1

start_service() {
    curl -f https://raw.githubusercontent.com/medvedicos/Install_MiSa/refs/heads/main/hosts \
         --output /etc/hosts
}
